# root/etc/s6-overlay/s6-rc.d/init-config-kopia/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p \
    /config/log

touch /config/rclone.conf

# permissions
chown -R abc:abc \
    /config \
    /tmp
