# root/etc/s6-overlay/s6-rc.d/init-config-kopia/run
#!/usr/bin/with-contenv bash

touch /config/rclone.conf

# make folders
mkdir -p \
    /config/log/

# permissions
chown -R abc:abc \
    /config
